# send-digest.cron
0 0 * * FRI NODE_ENV=production node index.js
